# packages/bsp/jethub/setupNetwork/setupwifi.service
[Unit]
Description=Setup wifi
After=wpa_supplicant-nl80211@wlan0.service
Requires=wpa_supplicant-nl80211@wlan0.service

[Service]
Type=oneshot
ExecStartPre=/usr/lib/thirdreality/setupwifi.sh check
ExecStart=/usr/lib/thirdreality/setupwifi.sh start
ExecStop=/usr/lib/thirdreality/setupwifi.sh stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target