# packages/bsp/jethub/w1/amlogicw1-bt.service
[Unit]
Description=Amlogic W155s1 Bluetooth Initialization Service
After=amlogicw1.service
Requires=amlogicw1.service

#After=wpa_supplicant-nl80211@wlan0.service
#Requires=wpa_supplicant-nl80211@wlan0.service

[Service]
Type=oneshot
ExecStart=/usr/lib/firmware/w1/w1_bt_start.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target