# packages/bsp/hubv3/linuxbox/hubv3-button.service
[Unit]
Description=HubV3 button Control Service
Before=basic.target
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/lib/armbian/hubv3-button-ctrl start
ExecStop=/lib/armbian/hubv3-button-ctrl stop
Restart=always
PIDFile=/run/button_check.pid
TimeoutStartSec=10

[Install]
WantedBy=basic.target
