# packages/bsp/thirdreality/trhubv3/hubv3-init.service
[Unit]
Description=HubV3 hardware initializer
Before=basic.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/lib/thirdreality/hubv3-init.sh
TimeoutStartSec=1min

[Install]
WantedBy=basic.target
