# packages/bsp/thirdreality/trhubv3/hubv3-monitor.service
[Unit]
Description=Hardware Monitor for LED and Button Control
Before=basic.target
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/usr/bin/env python3 /lib/thirdreality/hubv3-monitor.py daemon
Restart=always
User=root
Environment=PYTHONUNBUFFERED=1
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=hwmonitor

[Install]
WantedBy=multi-user.target
