# packages/bsp/thirdreality/common/hubv3-usb-sync.service
[Unit]
Description=Thirdreality USB auto sync Service
After=basic.target

[Service]
Type=simple
ExecStart=/usr/lib/thirdreality/hubv3-usb-sync.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hubv3-usb-sync
Restart=no
User=root

[Install]
WantedBy=multi-user.target
