# packages/bsp/thirdreality/common/hubv3-boot-monitor.service
[Unit]
Description=Hubv3 Boot wifi Monitor
After=network.target

[Service]
Type=simple
ExecStart=/lib/armbian/check-wifi-and-start-service.sh wlan0
Restart=on-failure

[Install]
WantedBy=multi-user.target

