# packages/bsp/thirdreality/common/hubv3-factory-reset.service
[Unit]
Description=Hubv3 Factory-reset
After=network.target

[Service]
Type=oneshot
ExecStart=/lib/armbian/factory-reset.sh
RemainAfterExit=true
TimeoutStartSec=180
StandardOutput=journal
StandardError=journal


[Install]
WantedBy=multi-user.target

