# packages/bsp/thirdreality/common/thirdreality-health-checker.service
[Unit]
Description=Thirdreality HA Health Check Service
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash /usr/lib/thirdreality/thirdreality-health-checker.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
